# services/battery-monitor.service
[Unit]
Description=HackberryPi CM5 Battery Monitor
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/battery-monitor.sh monitor
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

