# services/ghostpi-bot-ai.service
[Unit]
Description=GhostPi Next-Generation AI Monitoring Bot 2025
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ghostpi-bot-ai.sh monitor
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

