# services/self-healing-2025.service
[Unit]
Description=GhostPi Advanced Self-Healing System 2025
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ghostpi-self-heal-2025.sh monitor
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

