# services/self-healing-ai.service
[Unit]
Description=GhostPi Next-Generation AI Self-Healing System 2025
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ghostpi-self-heal-ai.sh monitor
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

